# deploy/supervisor.conf
[program:fx-assistant]
command=/opt/fx-assistant/venv/bin/gunicorn -c /opt/fx-assistant/deploy/gunicorn.conf.py wsgi:app
directory=/opt/fx-assistant/backend
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/fx-assistant/supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=FLASK_ENV="production",DATABASE_URL="sqlite:////var/data/fx-assistant/prod.db",SECRET_KEY="change-me-to-random-string",JWT_SECRET_KEY="change-me-to-another-random-string"
